# etc/profile.d/formilux.var
PKGROOT=/mnt/cdrom/pkg
LYNX_CFG=/etc/lynx.cfg
PATH=$HOME/bin:$PATH:/opt/bin:/opt/sbin
LD_LIBRARY_PATH=$HOME/lib${LD_LIBRARY_PATH:":$LD_LIBRARY_PATH"}
MANPATH=/usr/man:/usr/share/man
